# test/dep.service
[Unit]
Description=Test dep
[Service]
Type=oneshot
#ExecStart=echo second >> testfile
ExecStart=test/dep.sh
