# test/sv.service
[Unit]
Description=Test
Requires=dep.service
After=dep.service
[Service]
ExecStart=test/sv.sh
