# test/sv.service
[Unit]
Description=Test
Requires=dep.service
After=dep.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=test/sv.sh
