# lib/systemd/system/hi.service
[Unit]
Description=hello service
[Service]
ExecStart=echo hi
